# rootfs/etc/systemd/system/gitlfspull.timer
[Unit]
Description=git lfs pull on a schedule

[Timer]
OnBootSec=30
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target