# rootfs/etc/systemd/system/update-timer.timer
[Unit]
Description=Update gitlfspull timer

[Timer]
OnBootSec=2s
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target